# deploy/templates/rclone-archive.timer.j2
[Unit]
Description=Perform rclone archive

[Timer]
OnCalendar=03:00:00

[Install]
WantedBy=timers.target
