# deploy/templates/reverse-geocode.timer.j2
[Unit]
Description=Perform reverse geocode lookups for locations without metadata

[Timer]
OnCalendar=02:12:00

[Install]
WantedBy=timers.target
